# deploy/templates/supervisor_web.cfg
[supervisord]
logfile = %(here)s/logs/supervisor.log
loglevel = info
pidfile = %(here)s/run/supervisord.pid
directory = %(here)s

[unix_http_server]
file = %(here)s/run/supervisord.sock
chown = www-data:www-data
chmod = 0770

[supervisorctl]
serverurl = unix:///%(here)s/run/supervisord.sock

[rpcinterface:supervisor]
supervisor.rpcinterface_factory = supervisor.rpcinterface:make_main_rpcinterface



{% for proc_name in process_groups %}
[program:{{ proc_name }}]
command = {{ deploy_path }}/env/bin/paster serve --server-name={{ proc_name }}_%(process_num)02d {{ config_file }}
process_name = {{ proc_name }}_%(process_num)02d
numprocs={{ num_procs }}
autostart = true
startretries=10
autorestart=true
startsecs = 5
user = www-data
redirect_stderr = true
stdout_logfile = {{ log_path }}/{{ proc_name }}_%(process_num)02d.log
environment = PYTHONPATH='{{ python_path }}'

{% endfor %}

[program:notifier]
command = {{ deploy_path }}/env/bin/python {{ deploy_path }}/code/current/friendfund/tasks/notifier.py -f {{ deploy_path }}/code/current/config.ini
process_name = %(program_name)s
autostart = true
startretries=10
autorestart=true
startsecs = 5
user = www-data
redirect_stderr = true
stdout_logfile = {{ log_path }}/notification.log
environment = PYTHONPATH='{{ python_path }}'

[program:payment]
command = {{ deploy_path }}/env/bin/python {{ deploy_path }}/code/current/friendfund/tasks/payment_job.py -f {{ deploy_path }}/code/current/config.ini
process_name = %(program_name)s
autostart = true
startretries=3
autorestart=true
startsecs = 5
user = www-data
redirect_stderr = true
stdout_logfile = {{ log_path }}/payment_job.log
environment = PYTHONPATH='{{ python_path }}'

[program:cm]
command = {{ deploy_path }}/env/bin/python {{ deploy_path }}/code/current/friendfund/tasks/cache_refresher.py -f {{ deploy_path }}/code/current/config.ini
process_name = %(program_name)s
autostart = true
startretries=50
autorestart=true
startsecs = 5
user = www-data
redirect_stderr = true
stdout_logfile = {{ log_path }}/cache.log
environment = PYTHONPATH='{{ python_path }}'

[program:celery]
command = {{ deploy_path }}/env/bin/celery worker --app=friendfund.tasks.celerytasks -l info
process_name = %(program_name)s
autostart = true
startretries=50
autorestart=true
startsecs = 5
user = www-data
redirect_stderr = true
stdout_logfile = {{ log_path }}/celery.log
environment = PYTHONPATH='{{ python_path }}'


