# environments/live_miranda/supervisor_miranda.cfg
[supervisord]
logfile = /opt/www/ff_live/logs/supervisor.log
loglevel = info
pidfile = /opt/www/ff_live/run/supervisord.pid
directory = /opt/www/ff_live

[unix_http_server]
file = /opt/www/ff_live/run/supervisord.sock
chown = www-data:www-data
chmod = 0770

[supervisorctl]
serverurl = unix:////opt/www/ff_live/run/supervisord.sock

[rpcinterface:supervisor]
supervisor.rpcinterface_factory = supervisor.rpcinterface:make_main_rpcinterface

[group:miranda1]
programs=miranda_1,miranda_2,miranda_3
[group:miranda2]
programs=miranda_4,miranda_5,miranda_6


[program:miranda_1]
command = /server/pylons1.0/bin/paster serve /opt/www/ff_live/miranda_1_paste.ini
process_name = %(program_name)s
autostart = true
startretries=10
autorestart=true
startsecs = 5
user = www-data
redirect_stderr = true
stdout_logfile = /opt/www/ff_live/logs/miranda_paste_1.log
environment=PYTHONPATH=/opt/www/ff_live:/home/www-data/instances/p1_curr

[program:miranda_2]
command = /server/pylons1.0/bin/paster serve /opt/www/ff_live/miranda_2_paste.ini
process_name = %(program_name)s
autostart = true
startretries=10
autorestart=true
startsecs = 5
user = www-data
redirect_stderr = true
stdout_logfile = /opt/www/ff_live/logs/miranda_paste_2.log
environment=PYTHONPATH=/opt/www/ff_live:/home/www-data/instances/p2_curr

[program:miranda_3]
command = /server/pylons1.0/bin/paster serve /opt/www/ff_live/miranda_3_paste.ini
process_name = %(program_name)s
autostart = true
startretries=10
autorestart=true
startsecs = 5
user = www-data
redirect_stderr = true
stdout_logfile = /opt/www/ff_live/logs/miranda_paste_3.log
environment=PYTHONPATH=/opt/www/ff_live:/home/www-data/instances/p3_curr


[program:miranda_4]
command = /server/pylons1.0/bin/paster serve /opt/www/ff_live/miranda_4_paste.ini
process_name = %(program_name)s
autostart = true
startretries=10
autorestart=true
startsecs = 5
user = www-data
redirect_stderr = true
stdout_logfile = /opt/www/ff_live/logs/miranda_paste_4.log
environment=PYTHONPATH=/opt/www/ff_live:/home/www-data/instances/p4_curr

[program:miranda_5]
command = /server/pylons1.0/bin/paster serve /opt/www/ff_live/miranda_5_paste.ini
process_name = %(program_name)s
autostart = true
startretries=10
autorestart=true
startsecs = 5
user = www-data
redirect_stderr = true
stdout_logfile = /opt/www/ff_live/logs/miranda_paste_5.log
environment=PYTHONPATH=/opt/www/ff_live:/home/www-data/instances/p5_curr

[program:miranda_6]
command = /server/pylons1.0/bin/paster serve /opt/www/ff_live/miranda_6_paste.ini
process_name = %(program_name)s
autostart = true
startretries=10
autorestart=true
startsecs = 5
user = www-data
redirect_stderr = true
stdout_logfile = /opt/www/ff_live/logs/miranda_paste_6.log
environment=PYTHONPATH=/opt/www/ff_live:/home/www-data/instances/p6_curr


