# environments/live_ariel/supervisor_ariel.cfg
[supervisord]
logfile = /opt/www/ff_live/logs/supervisor.log
loglevel = info
pidfile = /opt/www/ff_live/run/supervisord.pid
directory = /opt/www/ff_live

[unix_http_server]
file = /opt/www/ff_live/run/supervisord.sock
chown = www-data:www-data
chmod = 0770

[supervisorctl]
serverurl = unix:////opt/www/ff_live/run/supervisord.sock

[rpcinterface:supervisor]
supervisor.rpcinterface_factory = supervisor.rpcinterface:make_main_rpcinterface

[group:ariel1]
programs=ariel_101,ariel_102,ariel_103
[group:ariel2]
programs=ariel_201,ariel_202,ariel_203
[group:jobs]
programs=notifier,payment


[program:ariel_101]
command = /server/pylons1.0/bin/paster serve /opt/www/ff_live/ariel_101_paste.ini
process_name = %(program_name)s
autostart = true
startretries=10
autorestart=true
startsecs = 5
user = www-data
redirect_stderr = true
stdout_logfile = /opt/www/ff_live/logs/ariel_paste_101.log
environment=PYTHONPATH=/opt/www/ff_live:/home/www-data/instances/p1_curr

[program:ariel_102]
command = /server/pylons1.0/bin/paster serve /opt/www/ff_live/ariel_102_paste.ini
process_name = %(program_name)s
autostart = true
startretries=10
autorestart=true
startsecs = 5
user = www-data
redirect_stderr = true
stdout_logfile = /opt/www/ff_live/logs/ariel_paste_102.log
environment=PYTHONPATH=/opt/www/ff_live:/home/www-data/instances/p1_curr

[program:ariel_103]
command = /server/pylons1.0/bin/paster serve /opt/www/ff_live/ariel_103_paste.ini
process_name = %(program_name)s
autostart = true
startretries=10
autorestart=true
startsecs = 5
user = www-data
redirect_stderr = true
stdout_logfile = /opt/www/ff_live/logs/ariel_paste_103.log
environment=PYTHONPATH=/opt/www/ff_live:/home/www-data/instances/p1_curr


[program:ariel_201]
command = /server/pylons1.0/bin/paster serve /opt/www/ff_live/ariel_201_paste.ini
process_name = %(program_name)s
autostart = true
startretries=10
autorestart=true
startsecs = 5
user = www-data
redirect_stderr = true
stdout_logfile = /opt/www/ff_live/logs/ariel_paste_201.log
environment=PYTHONPATH=/opt/www/ff_live:/home/www-data/instances/p2_curr

[program:ariel_202]
command = /server/pylons1.0/bin/paster serve /opt/www/ff_live/ariel_202_paste.ini
process_name = %(program_name)s
autostart = true
startretries=10
autorestart=true
startsecs = 5
user = www-data
redirect_stderr = true
stdout_logfile = /opt/www/ff_live/logs/ariel_paste_202.log
environment=PYTHONPATH=/opt/www/ff_live:/home/www-data/instances/p2_curr

[program:ariel_203]
command = /server/pylons1.0/bin/paster serve /opt/www/ff_live/ariel_203_paste.ini
process_name = %(program_name)s
autostart = true
startretries=10
autorestart=true
startsecs = 5
user = www-data
redirect_stderr = true
stdout_logfile = /opt/www/ff_live/logs/ariel_paste_203.log
environment=PYTHONPATH=/opt/www/ff_live:/home/www-data/instances/p2_curr



[program:notifier]
command = python /home/www-data/instances/notif_curr/friendfund/tasks/notifier.py -f /opt/www/ff_live/ariel_101_paste.ini
process_name = %(program_name)s
autostart = true
startretries=10
autorestart=true
startsecs = 5
user = www-data
redirect_stderr = true
stdout_logfile = /opt/www/ff_live/logs/notification.log
environment=PYTHONPATH=/opt/www/ff_live:/home/www-data/instances/notif_curr

[program:payment]
command = python /home/www-data/instances/notif_curr/friendfund/tasks/payment_job.py -f /opt/www/ff_live/ariel_101_paste.ini
process_name = %(program_name)s
autostart = true
startretries=10
autorestart=true
startsecs = 5
user = www-data
redirect_stderr = true
stdout_logfile = /opt/www/ff_live/logs/payment_job.log
environment=PYTHONPATH=/opt/www/ff_live:/home/www-data/instances/payment_curr







