# environments/live_ariel/supervisor_ariel.cfg
[supervisord]
logfile = /opt/www/ff_live/logs/supervisor.log
loglevel = info
pidfile = /opt/www/ff_live/run/supervisord.pid
directory = /opt/www/ff_live
environment=PYTHONPATH=/opt/www/ff_live:/home/www-data/deployment/ff_current

[unix_http_server]
file = /opt/www/ff_live/run/supervisord.sock
chown = www-data:www-data
chmod = 0770

[supervisorctl]
serverurl = unix:////opt/www/ff_live/run/supervisord.sock

[rpcinterface:supervisor]
supervisor.rpcinterface_factory = supervisor.rpcinterface:make_main_rpcinterface

[group:ariel]
programs=ariel_1,ariel_2

[group:notif]
programs=notifier


[program:ariel_1]
command = /server/pylons1.0/bin/paster serve /opt/www/ff_live/ariel_1_paste.ini
process_name = %(program_name)s
autostart = true
startretries=10
autorestart=true
startsecs = 5
user = www-data
redirect_stderr = true
stdout_logfile = /opt/www/ff_live/logs/ariel_paste_1.log
environment=PYTHONPATH=/opt/www/ff_live:/home/www-data/deployment/ff_current

[program:ariel_2]
command = /server/pylons1.0/bin/paster serve /opt/www/ff_live/ariel_2_paste.ini
process_name = %(program_name)s
autostart = true
startretries=10
autorestart=true
startsecs = 5
user = www-data
redirect_stderr = true
stdout_logfile = /opt/www/ff_live/logs/ariel_paste_2.log
environment=PYTHONPATH=/opt/www/ff_live:/home/www-data/deployment/ff_current


[program:notifier]
command = python /home/www-data/deployment/ff_current/friendfund/tasks/notifier.py -f /opt/www/ff_live/ariel_1_paste.ini
process_name = %(program_name)s
autostart = true
startretries=10
autorestart=true
startsecs = 5
user = www-data
redirect_stderr = true
stdout_logfile = /opt/www/ff_live/logs/notification.log
environment=PYTHONPATH=/opt/www/ff_live:/home/www-data/deployment/ff_current

