# environments/live_ariel/supervisor_ariel.cfg
[supervisord]
logfile = /opt/www/ff_live/logs/supervisor.log
loglevel = info
pidfile = /opt/www/ff_live/run/supervisord.pid
directory = /opt/www/ff_live

[unix_http_server]
file = /opt/www/ff_live/run/supervisord.sock
chown = www-data:www-data
chmod = 0770

[supervisorctl]
serverurl = unix:////opt/www/ff_live/run/supervisord.sock

[rpcinterface:supervisor]
supervisor.rpcinterface_factory = supervisor.rpcinterface:make_main_rpcinterface

[group:web]
programs=p1,p2
[group:jobs]
programs=notifier,payment,cm


[program:p1]
command = /server/pylons1.0/bin/paster serve --server-name=webp1_%(process_num)02d  /opt/www/ff_live/ariel_101_paste.ini
process_name = friendfund_%(program_name)s_%(process_num)02d
numprocs=3
autostart = true
startretries=10
autorestart=true
startsecs = 5
user = www-data
redirect_stderr = true
stdout_logfile = /opt/www/ff_live/logs/ariel_paste_p1.log
environment=PYTHONPATH=/opt/www/ff_live:/home/www-data/instances/p1_curr

[program:p2]
command = /server/pylons1.0/bin/paster serve --server-name=webp2_%(process_num)02d  /opt/www/ff_live/ariel_101_paste.ini
process_name = friendfund_%(program_name)s_%(process_num)02d
numprocs=3
autostart = true
startretries=10
autorestart=true
startsecs = 5
user = www-data
redirect_stderr = true
stdout_logfile = /opt/www/ff_live/logs/ariel_paste_p2.log
environment=PYTHONPATH=/opt/www/ff_live:/home/www-data/instances/p1_curr


[program:notifier]
command = python /home/www-data/instances/notif_curr/friendfund/tasks/notifier.py -f /opt/www/ff_live/ariel_101_paste.ini
process_name = %(program_name)s
autostart = true
startretries=10
autorestart=true
startsecs = 5
user = www-data
redirect_stderr = true
stdout_logfile = /opt/www/ff_live/logs/notification.log
environment=PYTHONPATH=/opt/www/ff_live:/home/www-data/instances/notif_curr

[program:payment]
command = python /home/www-data/instances/notif_curr/friendfund/tasks/payment_job.py -f /opt/www/ff_live/ariel_101_paste.ini
process_name = %(program_name)s
autostart = true
startretries=3
autorestart=true
startsecs = 5
user = www-data
redirect_stderr = true
stdout_logfile = /opt/www/ff_live/logs/payment_job.log
environment=PYTHONPATH=/opt/www/ff_live:/home/www-data/instances/notif_curr

[program:cm]
command = python /home/www-data/instances/notif_curr/friendfund/tasks/cache_refresher.py -f /opt/www/ff_live/ariel_101_paste.ini
process_name = %(program_name)s
autostart = true
startretries=50
autorestart=true
startsecs = 5
user = www-data
redirect_stderr = true
stdout_logfile = /opt/www/ff_live/logs/cache.log
environment=PYTHONPATH=/opt/www/ff_live:/home/www-data/instances/notif_curr
