# environments/live_diariobebe/supervisor_diariobebe.cfg
[supervisord]
logfile = /opt/www/diariobebe/logs/supervisor.log
loglevel = info
pidfile = /opt/www/diariobebe/run/supervisord.pid
directory = /opt/www/diariobebe


[unix_http_server]
file = /opt/www/diariobebe/run/supervisord.sock
chown = www-data:www-data
chmod = 0770

[supervisorctl]
serverurl = unix:////opt/www/diariobebe/run/supervisord.sock

[rpcinterface:supervisor]
supervisor.rpcinterface_factory = supervisor.rpcinterface:make_main_rpcinterface

[group:ariel1]
programs=ariel_1


[program:ariel_1]
command = /server/pylons1.0/bin/paster serve /opt/www/diariobebe/ariel_1_diariobebe.ini
process_name = %(program_name)s
autostart = true
startretries=10
autorestart=true
startsecs = 5
user = www-data
redirect_stderr = true
stdout_logfile = /opt/www/diariobebe/logs/ariel_diariobebe_1.log
environment=PYTHONPATH=/opt/www/diariobebe:/home/www-data/instances/diariobebe_curr:/opt/www/ff_live

